# ops/systemd/expresso-smoke-multirealm@.timer
[Unit]
Description=Periodic multi-realm smoke test (tenant=%i, every 10 min)

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true
Unit=expresso-smoke-multirealm@%i.service
RandomizedDelaySec=60s

[Install]
WantedBy=timers.target
