# ops/systemd/expresso-smoke-multirealm@.service
[Unit]
Description=Expresso multi-realm smoke test (tenant=%i)
Documentation=file:///opt/expresso/smoke-multirealm.sh
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/expresso/smoke-multirealm-%i.env
ExecStart=/bin/bash /opt/expresso/smoke-multirealm.sh
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
