# ops/backup/expresso-pg-backup.timer
[Unit]
Description=Expresso — run PostgreSQL backup daily 02:00 UTC

[Timer]
OnCalendar=*-*-* 02:00:00 UTC
Persistent=true
RandomizedDelaySec=300
Unit=expresso-pg-backup.service

[Install]
WantedBy=timers.target
